# scripts/robinhood.target
[Unit]
Description=Convenience target to manage all robinhood instances

[Install]
WantedBy=multi-user.target
