# scripts/robinhood.service.in
[Unit]
Description=Robinhood server
After=remote-fs.target
#only works if config file is unique

[Service]
Type=simple
KillMode=mixed
EnvironmentFile=-@CONFDIR@/sysconfig/robinhood
LimitNOFILE=8096
ExecStart=@SBINDIR@/robinhood $RBH_OPT
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
