# scripts/robinhood@.service.in
[Unit]
Description=Robinhood server for %I
After=local-fs.target remote-fs.target
Requires=local-fs.target
#to speedup the starting, remove the above two lines,
#  then uncomment and modify the following 'RequiresMountsFor' directive
#RequiresMountsFor=/path/to/scan

[Service]
Type=simple
KillMode=mixed
EnvironmentFile=-@CONFDIR@/sysconfig/robinhood
EnvironmentFile=-@CONFDIR@/sysconfig/robinhood.%I
LimitNOFILE=8096
ExecStart=@SBINDIR@/robinhood $RBH_OPT -f @CONFDIR@/robinhood.d/%I.conf
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
