# scripts/robinhood.service.in
[Unit]
Description=Robinhood server
#only works if config file is unique

[Service]
Type=simple
KillMode=mixed
EnvironmentFile=-@CONFDIR@/sysconfig/robinhood
ExecStart=@SBINDIR@/robinhood $RBH_OPT
ExecReload=/bin/kill -HUP $MAINPID
